# 3dprinter-octoprint/rootfs/etc/services.d/octoprint/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Add-on: 3dprinter-octoprint
# Runs OctoPrint
# ==============================================================================

# Lê a opção do config.yaml
SAFE_MODE=$(bashio::config 'safe_mode')

# Define os argumentos base
ARGS="--basedir /data/config/octoprint serve --iknowwhatimdoing --host 0.0.0.0 --port 80"

# Se o Safe Mode estiver ligado no HA, adiciona o flag
if bashio::var.true "${SAFE_MODE}"; then
    bashio::log.warning "MODO DE SEGURANÇA ATIVADO! Ignorando plugins de terceiros."
    ARGS="${ARGS} --safe"
fi

bashio::log.info "Starting OctoPrint..."

## Executa o OctoPrint com os argumentos dinâmicos
exec octoprint ${ARGS}
