# 3dprinter-octoprint/rootfs/etc/services.d/octoprint/run
#!/usr/bin/with-contenv bashio

SAFE_MODE=$(bashio::config 'safe_mode')
CAMERA_URL=$(bashio::config 'camera_url')
TRUSTED_PROXIES=$(bashio::config 'trusted_proxies')

ARGS="--basedir /data/.octoprint --host 0.0.0.0 --port 80 --iknowwhatimdoing"

if bashio::var.true "${SAFE_MODE}"; then
    bashio::log.warning "Starting in SAFE MODE! Third-party plugins are disabled."
    ARGS="${ARGS} --safe"
fi

bashio::log.info "Starting OctoPrint 1.11.6 (Python 3.13.5)..."

export OCTOPRINT_REVERSE_PROXY_TRUST_ADDRESSES="127.0.0.1,::1,172.30.32.1/23"

# Execute OctoPrint
exec octoprint serve ${ARGS}
