# 3dprinter-octoprint/rootfs/etc/services.d/octoprint/run
SAFE_MODE=$(bashio::config 'safe_mode')

ARGS="--basedir /data/config/octoprint serve --iknowwhatimdoing --host 0.0.0.0 --port 80"

if bashio::var.true "${SAFE_MODE}"; then
    bashio::log.warning "MODO DE SEGURANÇA ATIVADO! Ignorando plugins de terceiros."
    ARGS="${ARGS} --safe"
fi

bashio::log.info "Starting OctoPrint..."

exec octoprint ${ARGS}
