# 3dprinter-octoprint/rootfs/etc/services.d/octoprint/run
#!/usr/bin/with-contenv bashio

SAFE_MODE=$(bashio::config 'safe_mode')

# Reverter qualquer alteração de porta no config.yaml para a original 5000
if [ -f "/data/.octoprint/config.yaml" ]; then
    sed -i 's/port: [0-9]*/port: 5000/g' /data/.octoprint/config.yaml
fi

ARGS="--basedir /data/.octoprint --host 0.0.0.0 --port 5000 --iknowwhatimdoing"

if bashio::var.true "${SAFE_MODE}"; then
    bashio::log.warning "Starting in SAFE MODE!"
    ARGS="${ARGS} --safe"
fi

bashio::log.info "Starting OctoPrint on port 5000..."

exec octoprint serve ${ARGS}
